# deploy/systemd/workdir.service
[Unit]
Description=workdir — low-cost Firecracker microVM sandbox provider
Documentation=https://workdir.dev
After=network-online.target workdir-net.service
Wants=network-online.target workdir-net.service

[Service]
Type=simple
# Runs as root because the jailer (runtime.use_jailer=true, the default for the
# production provisioner) needs to chroot + drop each microVM to its own
# per-VM uid/gid. Each Firecracker process then runs UNPRIVILEGED and jailed —
# the daemon is root, the VMMs are not. With use_jailer=false you can instead
# run User=workdir (+ SupplementaryGroups=kvm, AmbientCapabilities=CAP_NET_ADMIN
# CAP_NET_BIND_SERVICE for the sandbox DNS proxy).
User=root
Environment=WORKDIR_CONFIG=/etc/workdir/config.toml
Environment=RUST_LOG=info,sandboxd=info
# --config is a global flag (before the subcommand); the unit passes it via the
# WORKDIR_CONFIG env var above, so `serve` takes no extra args.
ExecStart=/usr/local/bin/workdir serve
Restart=on-failure
RestartSec=2
LimitNOFILE=1048576
StateDirectory=workdir

[Install]
WantedBy=multi-user.target
